# scripts/launchd/ch.hydrosolutions.sapphire-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ch.hydrosolutions.sapphire-watchdog</string>
    <!-- Plan 147 Slice C: no --probe-token-path arg is added here
         DELIBERATELY. DEFAULT_PROBE_TOKEN_PATH
         (./secrets/health_probe_token, ops/watchdog.py) is relative and
         resolves against WorkingDirectory below to the correct HOST secret
         file — same convention already used for --slack-path (also
         unpassed, also default+WorkingDirectory). See
         docs/standards/cicd.md § Access-token pepper + probe-token
         rotation for the bootstrap runbook. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/sapphire/.local/bin/uv</string>
        <string>run</string>
        <string>--no-sync</string>
        <string>python</string>
        <string>-m</string>
        <string>sapphire_flow.ops.watchdog</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>WorkingDirectory</key>
    <string>/Users/sapphire/SAPPHIRE_flow</string>
    <key>StandardOutPath</key>
    <string>/Users/sapphire/Library/Logs/sapphire-watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sapphire/Library/Logs/sapphire-watchdog.log</string>
    <key>RunAtLoad</key>
    <false/>
    <key>UserName</key>
    <string>sapphire</string>
</dict>
</plist>
