# scripts/launchd/ch.hydrosolutions.sapphire-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ch.hydrosolutions.sapphire-watchdog</string>
    <!-- Plan 147 Slice C: --probe-token-path is passed EXPLICITLY here
         (defensive — do NOT rely on the argparse default silently matching
         WorkingDirectory). It resolves against WorkingDirectory below to the
         HOST secret file ./secrets/health_probe_token — an admin-scoped
         access-token raw key the watchdog sends as `Authorization: Bearer`
         on its /health/detail probes. See docs/standards/cicd.md
         § Access-token pepper + probe-token rotation for the bootstrap +
         rotation runbook. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/sapphire/.local/bin/uv</string>
        <string>run</string>
        <string>--no-sync</string>
        <string>python</string>
        <string>-m</string>
        <string>sapphire_flow.ops.watchdog</string>
        <string>--probe-token-path</string>
        <string>./secrets/health_probe_token</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>WorkingDirectory</key>
    <string>/Users/sapphire/SAPPHIRE_flow</string>
    <key>StandardOutPath</key>
    <string>/Users/sapphire/Library/Logs/sapphire-watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sapphire/Library/Logs/sapphire-watchdog.log</string>
    <key>RunAtLoad</key>
    <false/>
    <key>UserName</key>
    <string>sapphire</string>
</dict>
</plist>
