# scripts/launchd/ch.hydrosolutions.sapphire-recap-probe.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ch.hydrosolutions.sapphire-recap-probe</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sapphire/SAPPHIRE_flow/scripts/launchd/run-recap-probe.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>10800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>UserName</key>
    <string>sapphire</string>
    <key>StandardOutPath</key>
    <string>/Users/sapphire/Library/Logs/sapphire-recap-probe.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sapphire/Library/Logs/sapphire-recap-probe.launchd.log</string>
</dict>
</plist>
